# launchd/com.claude-rc-watchdog.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-rc-watchdog</string>

    <!--
      caffeinate -i wraps the watchdog so it holds a
      PreventUserIdleSystemSleep assertion for as long as the daemon
      runs. Without this, Mac enters Idle Sleep when all 6 claude panes
      go quiet (claude's own caffeinate -i -t 300 only fires per-Bash-
      tool, ~5 min lease), the relay websockets time out, and the
      watchdog only runs in recovery mode. With caffeinate -i, the
      Mac stays awake whenever the watchdog is loaded — the Mac can
      still sleep on lid close or scheduled sleep, just not on idle.

      The /bin/bash -c wrapper redirects stdout/stderr to the log file
      in APPEND mode (>>). Using StandardOutPath/StandardErrorPath
      directly would O_TRUNC the file every LaunchAgent boot, wiping
      the SIGTERM/modal-dismiss audit trail across restarts. With
      append mode, the watchdog's own "[ts] starting (...)" banner
      doubles as a boot marker.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>exec /usr/bin/caffeinate -i __REPO__/bin/claude-rc-watchdog >> /tmp/claude-rc-watchdog.log 2>&amp;1</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>LC_CTYPE</key>
        <string>en_US.UTF-8</string>
    </dict>
</dict>
</plist>
