# launchd/com.claude-rc-watchdog.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-rc-watchdog</string>

    <!--
      caffeinate -i wraps the watchdog so it holds a
      PreventUserIdleSystemSleep assertion for as long as the daemon
      runs. Without this, Mac enters Idle Sleep when all 6 claude panes
      go quiet (claude's own caffeinate -i -t 300 only fires per-Bash-
      tool, ~5 min lease), the relay websockets time out, and the
      watchdog only runs in recovery mode. With caffeinate -i, the
      Mac stays awake whenever the watchdog is loaded — the Mac can
      still sleep on lid close or scheduled sleep, just not on idle.
    -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-i</string>
        <string>__REPO__/bin/claude-rc-watchdog</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>LANG</key>
        <string>en_US.UTF-8</string>
        <key>LC_CTYPE</key>
        <string>en_US.UTF-8</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/claude-rc-watchdog.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/claude-rc-watchdog.log</string>
</dict>
</plist>
